# scripts/loophole-server.service
[Unit]
Description=Loophole Tunnel Server
Documentation=https://github.com/your-org/loophole
After=network.target

[Service]
Type=simple
User=loophole
Group=loophole
ExecStart=/usr/local/bin/tunnel-server --config /etc/loophole/server.toml
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/loophole

# Resource limits
LimitNOFILE=65535

# Environment
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
